# deploy/raspberry-control.service
[Unit]
Description=Raspberry Control Web Application
After=network.target raspberry-control-mqtt.service
Wants=raspberry-control-mqtt.service
Documentation=https://github.com/SirRellik/raspberry_control

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/raspberry_control
Environment="NODE_ENV=production"
Environment="PORT=3000"
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 3000
Restart=always
RestartSec=10

# Bezpečnostní nastavení
NoNewPrivileges=true
PrivateTmp=true

# Logování
StandardOutput=journal
StandardError=journal
SyslogIdentifier=raspberry-control

[Install]
WantedBy=multi-user.target
